# 100-sound-hdmi.rules
SUBSYSTEM=="drm", KERNEL=="card0", ACTION=="change", RUN+="/opt/sound.sh"
